# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid

; ── nginx (public-facing reverse proxy) ───────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=10

; ── Super Agent FastAPI (internal port 8001) ──────────────────────────────────
[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8001 --workers 1
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/uvicorn.log
stderr_logfile=/var/log/supervisor/uvicorn_err.log
priority=20

; ── code-server / VSCode (internal port 8080) ─────────────────────────────────
[program:code-server]
command=code-server --bind-addr 0.0.0.0:8080 --user-data-dir /workspace/.vscode --extensions-dir /workspace/.vscode-ext --auth password /workspace
directory=/workspace
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/code-server.log
stderr_logfile=/var/log/supervisor/code-server_err.log
environment=PASSWORD="%(ENV_CODE_SERVER_PASSWORD)s",HOME="/root"
priority=30
