# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; ── nginx (public-facing reverse proxy) ───────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Super Agent FastAPI (internal port 8001) ──────────────────────────────────
[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8001 --workers 1
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; ── code-server / VSCode (internal port 3001) ─────────────────────────────────
; Config (auth, password, bind-addr) is written to /root/.config/code-server/config.yaml
; by entrypoint.sh at startup — no env-var expansion needed here.
[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:3001 --auth password --disable-telemetry --disable-update-check /workspace
directory=/workspace
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",PASSWORD="%(ENV_UI_PASSWORD)s"
priority=30
