# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; ── nginx (public-facing reverse proxy) ───────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Super Agent FastAPI (internal port 8001) ──────────────────────────────────
[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8001 --workers 1
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
environment=HOME="/root"

; ── code-server / VSCode (internal port 3001) ─────────────────────────────────
;
; STABILITY DESIGN:
;   stopasgroup=true + killasgroup=true
;     When code-server is restarted, the ENTIRE process group is signalled —
;     this includes Claude Code CLI subprocesses, language servers, and terminal
;     shells. Without this, orphaned child processes accumulate over time,
;     exhaust memory/file descriptors, and eventually crash code-server.
;
;   startsecs=8
;     Supervisor waits 8 seconds before declaring the process "running".
;     Prevents rapid restart loops when code-server takes time to bind port 3001.
;
;   startretries=20
;     Retries up to 20 times before giving up — survives transient failures
;     without human intervention.
;
;   stopwaitsecs=15
;     Gives the process group 15 seconds to flush state before SIGKILL.
;
[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:3001 --auth password --disable-telemetry --disable-update-check /workspace
directory=/workspace
autostart=true
autorestart=true
startretries=20
startsecs=8
stopwaitsecs=15
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",PATH="/opt/flutter/bin:/opt/android-sdk/cmdline-tools/latest/bin:/opt/android-sdk/platform-tools:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",PASSWORD="%(ENV_UI_PASSWORD)s",PORT="3001",ANDROID_HOME="/opt/android-sdk",FLUTTER_HOME="/opt/flutter",JAVA_HOME="/usr/lib/jvm/java-17-openjdk-amd64"
priority=30
